# systemd/user/herdres-speech.service
[Unit]
Description=Herdres warm speech sidecar (STT + TTS) — issue #4

[Service]
Type=simple
EnvironmentFile=%h/.config/herdres/herdres.env
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HERDR_TELEGRAM_TOPICS_STATE=%h/.local/share/herdres/state.json
# Holds the parakeet (STT) + Kokoro (TTS) models warm so herdres command/sync get fast, warm
# transcription/synthesis over the Unix socket instead of a cold per-call model load. Opt-in:
# this unit is installed but NOT enabled by default. Enable after `herdres speech install`.
ExecStart=%h/.local/bin/herdres-speech
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target
