# systemd/user/herdres.timer
[Unit]
Description=Run Herdres Tendwire source Telegram sync

[Timer]
OnBootSec=30s
OnUnitActiveSec=10s
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
